# packaging/postinst
#!/bin/sh
# postinst -- configure proxmox-redfish. No venv/pip/network: the app runs on the
# system python3 with apt-provided requests/cryptography + bundled pure-Python deps.
set -e

ETC=/etc/proxmox-redfish
STATE=/var/lib/proxmox-redfish

case "$1" in
    configure)
        # Runtime directories.
        mkdir -p "$STATE/secureboot" "$STATE/varstores"
        chmod 0750 "$STATE" "$STATE/secureboot"
        mkdir -p "$ETC"

        # Ship the SecureBoot profile catalog example if no live one exists.
        if [ ! -f "$ETC/secureboot_profiles.json" ] && \
           [ -f /opt/proxmox-redfish/config/secureboot_profiles.json.example ]; then
            cp /opt/proxmox-redfish/config/secureboot_profiles.json.example "$ETC/secureboot_profiles.json"
        fi

        # Generate a self-signed TLS cert for the Redfish endpoint if none present.
        if [ ! -f "$ETC/server.crt" ] || [ ! -f "$ETC/server.key" ]; then
            CN="$(hostname -f 2>/dev/null || hostname)"
            openssl req -x509 -newkey rsa:4096 -nodes -days 825 \
                -keyout "$ETC/server.key" -out "$ETC/server.crt" \
                -subj "/CN=${CN}" >/dev/null 2>&1 || true
            chmod 0640 "$ETC/server.key" 2>/dev/null || true
        fi
        chmod 0640 "$ETC/params.env" 2>/dev/null || true

        # Sanity-check that the app imports on the system interpreter.
        if ! PYTHONPATH=/opt/proxmox-redfish/src:/opt/proxmox-redfish/vendor \
                /usr/bin/python3 -c "import proxmox_redfish.proxmox_redfish" >/dev/null 2>&1; then
            echo "proxmox-redfish: WARNING - import check failed; verify python3-requests is installed." >&2
        fi

        # Enable (but do not start -- the admin must set credentials first).
        if [ -d /run/systemd/system ]; then
            systemctl daemon-reload || true
            systemctl enable proxmox-redfish.service >/dev/null 2>&1 || true
        fi

        echo "proxmox-redfish installed."
        echo "  1. Edit ${ETC}/params.env (set PROXMOX_HOST/USER/PASSWORD)."
        echo "  2. systemctl start proxmox-redfish"
        echo "  Endpoint: https://<host>:8443/redfish/v1"
        ;;
    abort-upgrade|abort-remove|abort-deconfigure)
        ;;
esac

exit 0
