# packaging/proxmox-redfish.service
[Unit]
Description=Proxmox Redfish Daemon
Documentation=https://github.com/jaybinks/proxmox-redfish
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# Root is required: SecureBoot enrollment writes the VM efidisk LVM volume (dd).
User=root
EnvironmentFile=/etc/proxmox-redfish/params.env
# System python3 + bundled pure-Python deps (proxmoxer, requests-toolbelt). The
# apt-provided python3-requests / python3-cryptography are on sys.path already.
Environment=PYTHONPATH=/opt/proxmox-redfish/src:/opt/proxmox-redfish/vendor
ExecStart=/usr/bin/python3 /opt/proxmox-redfish/src/proxmox_redfish/proxmox_redfish.py --port 8443
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
